# ansible/roles/bridge_host/templates/bridge_v2.service.j2
[Unit]
Description=Neohabitat Habitat protocol bridge (bridge_v2)
Documentation=https://github.com/frandallfarmer/neohabitat
# docker.service brings up mongo / habiproxy / otel-collector that the
# bridge dials on startup. We don't *require* them — bridge_v2 has its
# own elko-reachability backoff — but ordering avoids a noisy probe.
After=network-online.target docker.service
Wants=network-online.target

[Service]
# Type=notify so the new tableflip child can call sd_notify(MAINPID=...)
# and have systemd update its main-PID tracking before the parent
# exits. With Type=simple, KillMode=control-group nukes the entire
# cgroup when the parent dies, killing the tableflip child along with
# it — service goes "Deactivated successfully" 1s after every reload.
# NotifyAccess=all is required because the tableflip child isn't
# systemd's direct descendant (the parent forked it), so the default
# NotifyAccess=main rejects its sd_notify.
Type=notify
NotifyAccess=all
User={{ bridge_host_user }}
Group={{ bridge_host_group }}
ExecStart={{ bridge_host_binary_path }} \
  --graceful \
  --qlink \
  --otel.enabled \
{% for addr in bridge_host_listen_addrs %}
  --listen={{ addr }} \
{% endfor %}
  --elko.host={{ bridge_host_elko_host }} \
  --mongo.host={{ bridge_host_mongo_url }} \
  --mongo.db={{ bridge_host_mongo_db }} \
  --mongo.collection={{ bridge_host_mongo_collection }} \
  --context={{ bridge_host_initial_context }} \
  --rate={{ bridge_host_rate }} \
  --log.level={{ bridge_host_log_level }}
# `systemctl reload bridge_v2` → SIGHUP → tableflip graceful upgrade.
# The new binary on disk (placed by the deploy job at
# {{ bridge_host_binary_path }}) takes over via fork+exec; sessions
# survive via TCP_REPAIR. This is the whole reason we run on the host
# instead of in docker.
ExecReload=/bin/kill -HUP $MAINPID
# tableflip's child re-exec returns exit 0 from the parent. Don't
# treat that as a crash. Real crashes are non-zero — those we restart.
Restart=on-failure
RestartSec=5s
# TCP_REPAIR setsockopt requires CAP_NET_ADMIN. With User=root this is
# already granted; keep the explicit declaration so a future move to a
# non-root user only needs to flip User=.
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN
# OTLP / OTel SDK config — same values the in-container service used.
Environment=NEOHABITAT_ORIGINAL_HATCHERY={{ bridge_host_original_hatchery | ternary('true', 'false') }}
Environment=OTEL_EXPORTER_OTLP_ENDPOINT={{ bridge_host_otel_endpoint }}
Environment=OTEL_SERVICE_NAME=bridge_v2
Environment=OTEL_METRIC_EXPORT_INTERVAL=10000

[Install]
WantedBy=multi-user.target
